# templates/etc/systemd/system/_mountpoint_.automount.j2
[Unit]
Description=Automount BTRFS filesystem on {{ mount_entry.device }}

[Automount]
Where=/mnt/btrfs/{{ mount_entry.device | basename }}
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
